# passwall-init-v2
#!/bin/sh /etc/rc.common

start(){

sh /root/install-passwall-panel.sh

}

reload(){

sh /root/up.sh

}
